# svc/apcupsd.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='apcupsd'>

  <service 
	name='system/apcupsd' 
	type='service' 
	version='0'>
	
	<create_default_instance enabled="true"/>

    	<!--
       	 	First of all, if the config file is not present,
       	 	then we needn't bother with anything else.  
    	-->

    	<dependency 
		name="apcupsd-config-file" 
		grouping="require_all" 
		restart_on="none" 
		type="path">
	      <service_fmri value="file:///etc/opt/apcupsd/apcupsd.conf"/>
    	</dependency>

	<!--
		Since apcupsd may be providing a home directory service,
		it is as well that we ensure that the file-systems are
		all mounted before it is started. This is not essential
		but in general it is a good thing and doesn't really hurt. 
	-->
    
	<dependency 
		name='fs-local' 
		grouping='require_all' 
		restart_on='none' 
		type='service'>
		<service_fmri value='svc:/system/filesystem/local'/>
	</dependency>
    
	<!--
      		now we have the start stop and refresh methods
    	-->

    
	<exec_method 
		name='start' 
		type='method' 
		exec='/etc/init.d/apcupsd start' 
		timeout_seconds='30'>
	</exec_method>
    
	<exec_method 
		name='stop' 
		type='method' 
		exec='/etc/init.d/apcupsd stop' 
		timeout_seconds='30'>
	</exec_method>
	
	<exec_method 
		name='refresh' 
		type='method' 
		exec='/etc/init.d/apcupsd stop ; /etc/init.d/apcupsd start' 
		timeout_seconds='60'>
	</exec_method>	
	
	<exec_method 
		name='restart' 
		type='method' 
		exec='/etc/init.d/apcupsd stop ; /etc/init.d/apcupsd start' 
		timeout_seconds='60'>
	</exec_method>	

	<template>
		<common_name>
          		<loctext xml:lang='C'>
				APC UPS Daemon
			</loctext>
        	</common_name>
	</template>

  </service>
</service_bundle>
